# systemd/pyhss_hss.service
[Unit]
Description=PyHSS HSS Service
PartOf=pyhss.service


[Service]
User=root
WorkingDirectory=/root/pyhss/services/
ExecStart=/root/pyhss/project_venv/bin/python hssService.py
Restart=always

[Install]
WantedBy=pyhss.service