# systemd/pyhss_geored.service
[Unit]
Description=PyHSS Geored Service
PartOf=pyhss.service


[Service]
User=root
WorkingDirectory=/root/pyhss/services/
ExecStart=/root/pyhss/project_venv/bin/python3.11 georedService.py
Restart=always

[Install]
WantedBy=pyhss.service