# systemd/pyhss_database.service
[Unit]
Description=PyHSS Database Service
PartOf=pyhss.service


[Service]
User=root
WorkingDirectory=/root/pyhss/services/
ExecStart=/root/pyhss/project_venv/bin/python3.11 databaseService.py
Restart=always

[Install]
WantedBy=pyhss.service